# tee-supplicant/tee-supplicant@.service.in
# SPDX-License-Identifier: BSD-2-Clause
[Unit]
Description=TEE Supplicant on %i
DefaultDependencies=no
After=dev-%i.device
Wants=dev-%i.device
Conflicts=shutdown.target
Before=tpm2.target sysinit.target shutdown.target

[Service]
Type=forking
User=@CFG_TEE_SUPPL_USER@
Group=@CFG_TEE_SUPPL_GROUP@
EnvironmentFile=-/@CMAKE_INSTALL_SYSCONFDIR@/default/tee-supplicant
ExecStart=@CMAKE_INSTALL_USR_SBINDIR@/tee-supplicant -d /dev/%i -l /lib/firmware/optee $OPTARGS
# start ftpm driver after tee-supplicant
ExecStartPost=-/usr/sbin/rmmod tpm_ftpm_tee ; /usr/sbin/modprobe tpm_ftpm_tee
KillSignal=SIGTERM
# Workaround for fTPM TA: stop kernel module before tee-supplicant
ExecStop=-/usr/sbin/modprobe -v -r tpm_ftpm_tee
